# root/etc/uci-defaults/99-campus-network-login
#!/bin/sh

chmod +x /usr/libexec/campus_login 2>&1 &
exit 0