# root/etc/uci-defaults/luci-campus-network-login
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@campus_network[-1]
	add ucitrack campus_network
	set ucitrack.@campus_network[-1].init=campus_network
	commit ucitrack
EOF

exit 0